# kobo/10-on-mount.rules
ACTION=="add", KERNEL=="sd*", ATTRS{vendor}=="Kobo    ", ATTRS{model}=="eReader-4.6.9995", RUN+="/kobo-sync.sh"
